# local.d/bnwifi.start
#!/bin/bash

count=0;

while [ $count -ne 4 ] ; do
    if [[ `wpa_cli scan && wpa_cli scan_results | egrep -io bnwifi` ]] ; then
        /bin/bash /home/anthony/Documents/Bash/BNWifi.sh 
	count=4;
	exit;
    else
        sleep 5;
    fi
    ((count++))
done

# Contents of: /home/anthony/Documents/Bash/BNWifi.sh
##!/bin/bash

#wpa_cli flush
#wpa_cli add_network
#wpa_cli set_network 0 ssid "BNWifi"
#wpa_cli set_network 0 key_mgmt NONE
#wpa_cli enable_network 0
#wpa_cli quit
